# roles/services/templates/flatpak-automatic.timer.j2
# {{ ansible_managed }}

[Unit]
Description=flatpak Automatic Update Trigger
Documentation=man:flatpak(1)

[Timer]
OnBootSec=3m
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
